# Tools/FightingVision2020.service
[Unit]

Description=...

[Service]

Execstart=/bin/FightingVision2020.sh
Restart=always
Type=simple

[Install]

WantedBy = multi-user.target